# src/data/remote/profiles/qilin.json
{
  "added_date": null,
  "altname": null,
  "description": "Qilin ransomware was first observed in July of 2022. Qilin Ransomware is written in Golang and supports multiple encryption modes; all of which are controlled by the operator. Qilin actors practice double extortion – demanding payment for a decryptor, as well as for the non-release of stolen data.",
  "locations": [
    {
      "available": false,
      "enabled": false,
      "fqdn": "kbsqoivihgdmwczmxkbovk7ss2dcynitwhhfu5yw725dboqo5kthfaad.onion",
      "slug": "http://kbsqoivihgdmwczmxkbovk7ss2dcynitwhhfu5yw725dboqo5kthfaad.onion",
      "title": "DDOS Protection",
      "type": "DLS"
    },
    {
      "available": true,
      "enabled": true,
      "fqdn": "ijzn3sicrcy7guixkzjkib4ukbiilwc3xhnmby4mcbccnsd7j2rekvqd.onion",
      "slug": "http://ijzn3sicrcy7guixkzjkib4ukbiilwc3xhnmby4mcbccnsd7j2rekvqd.onion",
      "title": "Qilin blog",
      "type": "DLS"
    },
    {
      "available": false,
      "enabled": false,
      "fqdn": "ji57fr53anp7wb44tbbnp72qcgbhqywy4jmbncawdcrejj5amuvh3zqd.onion",
      "slug": "http://ji57fr53anp7wb44tbbnp72qcgbhqywy4jmbncawdcrejj5amuvh3zqd.onion/site/login",
      "title": "Sign In",
      "type": "Admin"
    },
    {
      "available": false,
      "enabled": false,
      "fqdn": "ozsxj4hwxub7gio347ac7tyqqozvfioty37skqilzo2oqfs4cw2mgtyd.onion",
      "slug": "http://ozsxj4hwxub7gio347ac7tyqqozvfioty37skqilzo2oqfs4cw2mgtyd.onion",
      "title": "Qilin",
      "type": "DLS"
    }
  ],
  "name": "qilin",
  "tools": [
    {
      "CredentialTheft": [
        "Mimikatz"
      ],
      "DefenseEvasion": [
        "EDRSandBlast",
        "PCHunter",
        "PowerTool",
        "Toshiba power management driver (BYOVD)",
        "Updater for Carbon Black’s Cloud Sensor AV (upd.exe)",
        "YDArk",
        "Zemana Anti-Rootkit driver"
      ],
      "DiscoveryEnum": [
        "Nmap",
        "Nping"
      ],
      "Exfiltration": [
        "EasyUpload.io",
        "MEGA"
      ],
      "LOLBAS": [
        "PowerShell",
        "PsExec",
        "WinRM",
        "fsutil"
      ],
      "Networking": [
        "Proxychains"
      ],
      "Offsec": [
        "Cobalt Strike",
        "Evilginx",
        "Kali Linux",
        "NetExec",
        "SystemBC",
        "Tofsee"
      ],
      "RMM-Tools": [
        "NetSupport",
        "ScreenConnect"
      ]
    }
  ],
  "ttps": [
    {
      "tactic_id": "TA0001",
      "tactic_name": "Initial Access",
      "techniques": [
        {
          "technique_details": "Compromised credentials used to authenticate via VPN and RDP. Qilin notably steals credentials from Chrome browsers on compromised endpoints before encryption.",
          "technique_id": "T1078",
          "technique_name": "Valid Accounts"
        },
        {
          "technique_details": "Exploitation of vulnerabilities in VPN appliances and remote management tools to gain initial access.",
          "technique_id": "T1190",
          "technique_name": "Exploit Public-Facing Application"
        },
        {
          "technique_details": "",
          "technique_id": "T1566",
          "technique_name": "Phishing"
        },
        {
          "technique_details": "Qilin affiliates spear-phished MSP administrators via ScreenConnect to gain access to downstream customers. Phishing campaigns using credential-harvesting pages targeting VPN and SSO portals.",
          "technique_id": "T1566.003",
          "technique_name": "Phishing: Spearphishing via Service"
        }
      ]
    },
    {
      "tactic_id": "TA0002",
      "tactic_name": "Execution",
      "techniques": [
        {
          "technique_details": "PowerShell scripts used for payload execution, lateral movement, and post-exploitation tasks.",
          "technique_id": "T1059.001",
          "technique_name": "Command and Scripting Interpreter: PowerShell"
        },
        {
          "technique_details": "Shell scripts used to execute the Linux/ESXi variant of the Qilin ransomware payload.",
          "technique_id": "T1059.004",
          "technique_name": "Command and Scripting Interpreter: Unix Shell"
        },
        {
          "technique_details": "",
          "technique_id": "T1569",
          "technique_name": "System Services"
        },
        {
          "technique_details": "",
          "technique_id": "T1569.002",
          "technique_name": "System Services: Service Execution"
        }
      ]
    },
    {
      "tactic_id": "TA0003",
      "tactic_name": "Persistence",
      "techniques": [
        {
          "technique_details": "",
          "technique_id": "T1037",
          "technique_name": "Boot or Logon Initialization Scripts"
        },
        {
          "technique_details": "",
          "technique_id": "T1053",
          "technique_name": "Scheduled Task/Job"
        },
        {
          "technique_details": "Scheduled tasks created to maintain persistence and execute post-exploitation scripts.",
          "technique_id": "T1053.005",
          "technique_name": "Scheduled Task/Job: Scheduled Task"
        },
        {
          "technique_details": "",
          "technique_id": "T1098.004",
          "technique_name": "Account Manipulation: SSH Authorized Keys"
        },
        {
          "technique_details": "",
          "technique_id": "T1136",
          "technique_name": "Create Account"
        },
        {
          "technique_details": "",
          "technique_id": "T1547",
          "technique_name": "Boot or Logon Autostart Execution"
        }
      ]
    },
    {
      "tactic_id": "TA0004",
      "tactic_name": "Privilege Escalation",
      "techniques": [
        {
          "technique_details": "",
          "technique_id": "T1068",
          "technique_name": "Exploitation for Privilege Escalation"
        }
      ]
    },
    {
      "tactic_id": "TA0005",
      "tactic_name": "Defense Evasion",
      "techniques": [
        {
          "technique_details": "Qilin.B variant (October 2024) introduced enhanced obfuscation and anti-analysis techniques compared to earlier versions.",
          "technique_id": "T1027",
          "technique_name": "Obfuscated Files or Information"
        },
        {
          "technique_details": "",
          "technique_id": "T1036.001",
          "technique_name": "Masquerading: Invalid Code Signature"
        },
        {
          "technique_details": "",
          "technique_id": "T1134.004",
          "technique_name": "Access Token Manipulation: Parent PID Spoofing"
        },
        {
          "technique_details": "",
          "technique_id": "T1211",
          "technique_name": "Exploitation for Defense Evasion"
        },
        {
          "technique_details": "",
          "technique_id": "T1480",
          "technique_name": "Execution Guardrails"
        },
        {
          "technique_details": "",
          "technique_id": "T1497.001",
          "technique_name": "Virtualization/Sandbox Evasion: System Checks"
        },
        {
          "technique_details": "",
          "technique_id": "T1553.002",
          "technique_name": "Subvert Trust Controls: Code Signing"
        },
        {
          "technique_details": "Security tools and EDR solutions disabled via batch scripts and Group Policy modifications.",
          "technique_id": "T1562.001",
          "technique_name": "Disable or Modify Tools"
        },
        {
          "technique_details": "",
          "technique_id": "T1562.004",
          "technique_name": "Impair Defenses: Disable or Modify System Firewall"
        },
        {
          "technique_details": "",
          "technique_id": "T1564",
          "technique_name": "Hidden Artifacts"
        },
        {
          "technique_details": "",
          "technique_id": "T1564.003",
          "technique_name": "Hidden Artifacts: Hidden Window"
        }
      ]
    },
    {
      "tactic_id": "TA0006",
      "tactic_name": "Credential Access",
      "techniques": [
        {
          "technique_details": "LSASS memory dumped to harvest credentials for lateral movement and privilege escalation.",
          "technique_id": "T1003.001",
          "technique_name": "OS Credential Dumping: LSASS Memory"
        },
        {
          "technique_details": "",
          "technique_id": "T1040",
          "technique_name": "Network Sniffing"
        },
        {
          "technique_details": "",
          "technique_id": "T1110.002",
          "technique_name": "Brute Force: Password Cracking"
        },
        {
          "technique_details": "Qilin notably deployed a custom script to steal credentials stored in Google Chrome browsers on all domain-joined machines via Group Policy — a distinctive TTP first observed in the Synnovis NHS attack (June 2024).",
          "technique_id": "T1555.003",
          "technique_name": "Credentials from Web Browsers"
        }
      ]
    },
    {
      "tactic_id": "TA0007",
      "tactic_name": "Discovery",
      "techniques": [
        {
          "technique_details": "",
          "technique_id": "T1012",
          "technique_name": "Query Registry"
        },
        {
          "technique_details": "Network scanning used to enumerate hosts and services within the victim environment.",
          "technique_id": "T1046",
          "technique_name": "Network Service Discovery"
        },
        {
          "technique_details": "System information collected to tailor the attack and identify high-value targets.",
          "technique_id": "T1082",
          "technique_name": "System Information Discovery"
        },
        {
          "technique_details": "",
          "technique_id": "T1614",
          "technique_name": "System Location Discovery"
        }
      ]
    },
    {
      "tactic_id": "TA0008",
      "tactic_name": "Lateral Movement",
      "techniques": [
        {
          "technique_details": "",
          "technique_id": "T1021",
          "technique_name": "Remote Services"
        },
        {
          "technique_details": "RDP used for lateral movement across victim networks.",
          "technique_id": "T1021.001",
          "technique_name": "Remote Services: Remote Desktop Protocol"
        },
        {
          "technique_details": "SMB used to propagate payloads and move laterally within the network.",
          "technique_id": "T1021.002",
          "technique_name": "Remote Services: SMB/Windows Admin Shares"
        },
        {
          "technique_details": "",
          "technique_id": "T1021.004",
          "technique_name": "Remote Services: SSH"
        },
        {
          "technique_details": "",
          "technique_id": "T1570",
          "technique_name": "Lateral Tool Transfer"
        }
      ]
    },
    {
      "tactic_id": "TA0009",
      "tactic_name": "Collection",
      "techniques": [
        {
          "technique_details": "Data compressed and archived prior to exfiltration.",
          "technique_id": "T1560.001",
          "technique_name": "Archive Collected Data: Archive via Utility"
        },
        {
          "technique_details": "",
          "technique_id": "T1602.002",
          "technique_name": "Network Device Configuration Dump"
        }
      ]
    },
    {
      "tactic_id": "TA0010",
      "tactic_name": "Exfiltration",
      "techniques": [
        {
          "technique_details": "",
          "technique_id": "T1011",
          "technique_name": "Exfiltration Over Other Network Medium"
        },
        {
          "technique_details": "",
          "technique_id": "T1011.001",
          "technique_name": "Exfiltration Over Other Network Medium: Exfiltration Over Bluetooth"
        },
        {
          "technique_details": "",
          "technique_id": "T1048.003",
          "technique_name": "Exfiltration Over Alternative Protocol: Unencrypted Non-C2 Protocol"
        },
        {
          "technique_details": "Data exfiltrated to cloud storage and actor-controlled servers for double extortion via the Qilin leak site.",
          "technique_id": "T1567.002",
          "technique_name": "Exfiltration Over Web Service: Exfiltration to Cloud Storage"
        }
      ]
    },
    {
      "tactic_id": "TA0011",
      "tactic_name": "Command and Control",
      "techniques": [
        {
          "technique_details": "",
          "technique_id": "T1001",
          "technique_name": "Data Obfuscation"
        },
        {
          "technique_details": "",
          "technique_id": "T1001.001",
          "technique_name": "Data Obfuscation: Junk Data"
        },
        {
          "technique_details": "C2 communication over HTTPS; Tor used for victim portals.",
          "technique_id": "T1071.001",
          "technique_name": "Application Layer Protocol: Web Protocols"
        },
        {
          "technique_details": "",
          "technique_id": "T1572",
          "technique_name": "Protocol Tunneling"
        }
      ]
    },
    {
      "tactic_id": "TA0040",
      "tactic_name": "Impact",
      "techniques": [
        {
          "technique_details": "Qilin (formerly Agenda) uses ChaCha20 for file encryption with RSA-4096 for key protection. Originally written in Go, later rewritten in Rust (Qilin.B). Targets Windows, Linux, and VMware ESXi. Notable attack: Synnovis NHS pathology services (June 2024) causing 10,000+ NHS appointment cancellations. Active RaaS with high affiliate revenue share.",
          "technique_id": "T1486",
          "technique_name": "Data Encrypted for Impact"
        },
        {
          "technique_details": "Shadow copies deleted; ESXi snapshots removed to prevent VM recovery.",
          "technique_id": "T1490",
          "technique_name": "Inhibit System Recovery"
        },
        {
          "technique_details": "",
          "technique_id": "T1561",
          "technique_name": "Disk Wipe"
        },
        {
          "technique_details": "",
          "technique_id": "T1561.001",
          "technique_name": "Disk Wipe: Disk Content Wipe"
        }
      ]
    },
    {
      "tactic_id": "TA0042",
      "tactic_name": "Resource Development",
      "techniques": [
        {
          "technique_details": "",
          "technique_id": "T1587.001",
          "technique_name": "Develop Capabilities: Malware"
        }
      ]
    },
    {
      "tactic_id": "TA0043",
      "tactic_name": "Reconnaissance",
      "techniques": [
        {
          "technique_details": "",
          "technique_id": "T1590.004",
          "technique_name": "Gather Victim Network Information: Network Topology"
        }
      ]
    }
  ],
  "url": "https://www.ransomware.live/group/qilin"
}